# templates/cron/cron.timer.epp
<%- | String $on_calendar = 'daily',
      Optional[Variant[String,Boolean[false]]] $randomized_delay_sec = undef,
| -%>
[Timer]
OnCalendar=<%= $on_calendar %>
<% if $randomized_delay_sec { -%>
RandomizedDelaySec=<%= $randomized_delay_sec %>
<% } -%>

[Install]
WantedBy=timers.target
