# templates/cron/cron.timer.epp
<%- | String $on_calendar = 'daily',
      Optional[String] $randomize_delay_sec = undef,
| -%>
[Timer]
OnCalendar=<%= $on_calendar %>
<% if $randomize_delay_sec { -%>
RandomizedDelaySec=<%= $randomize_delay_sec %>
<% } -%>

[Install]
WantedBy=timers.target
